# macos/com.rac2.dotfiles-watch.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.rac2.dotfiles-watch</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>HOME_PLACEHOLDER/dotfiles/scripts/dotfiles-watch.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>HOME_PLACEHOLDER</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>HOME_PLACEHOLDER/.claude/dotfiles-watch.log</string>

  <key>StandardErrorPath</key>
  <string>HOME_PLACEHOLDER/.claude/dotfiles-watch.log</string>
</dict>
</plist>
